# deploy/systemd/smart-home-control.service
[Unit]
Description=Smart Home Control Gunicorn Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple

# Template values to replace on a new host:
# __APP_USER__
# __APP_GROUP__
# __BACKEND_DIR__
# __WEB_PORT__
# __GUNICORN_WORKERS__
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__BACKEND_DIR__
EnvironmentFile=__BACKEND_DIR__/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__BACKEND_DIR__/.venv/bin/gunicorn --workers __GUNICORN_WORKERS__ --bind 127.0.0.1:__WEB_PORT__ --access-logfile - --error-logfile - wsgi:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
